# buttonProb/buttonProb.rules
ACTION=="add", ATTRS{idVendor}=="04c5", ATTRS{idProduct}=="128d", RUN+="/root/sane-button/buttonProb.docrun"
