# tif2pdf/tif2pdf.socket
[Unit]
Description=listen pipe

[Socket]
ListenFIFO=/tmp/tif2pdf.pipe
RemoveOnStop=on

[Install]
WantedBy=tif2pdf.service

